# deploy/driftnote-backup-failure.service
[Unit]
Description=Driftnote backup failure alert

[Service]
Type=oneshot
EnvironmentFile=/etc/driftnote/driftnote.env
ExecStart=/usr/local/lib/driftnote/scripts/alert-email.py "Driftnote backup failed" "See journalctl -u driftnote-backup.service for details."
