# deploy/driftnote-backup.service
# Install: copy to ~/.config/systemd/user/driftnote-backup.service.
[Unit]
Description=Driftnote monthly backup
OnFailure=driftnote-backup-failure.service

[Service]
Type=oneshot
EnvironmentFile=%h/driftnote/driftnote.env
ExecStart=%h/.local/lib/driftnote/scripts/backup.sh
